# bin/crontab/backup-notes
#!/usr/bin/env sh

set -x
set -e

DIRS="$NOTE_DIR $YOUTUBE_DIR"
NOTE_DIR="$HOME/Documents/github/notes"
YOUTUBE_DIR="$HOME/Documents/github/youtube"

backup_with_git() {
	dir="$1"

	if [ ! -d "$dir" ]; then
		echo "Error: Directory '$dir' not found or inaccessible"
		return 1
	fi

	if git -C "$dir" cherry -v origin/main 2>/dev/null | grep -q .; then
		echo "Pushing existing commits for $dir"
		git -C "$dir" push origin main
	fi

	if git -C "$dir" status | grep -q "working tree clean"; then
		echo "Info: Nothing new to add to $dir"
		echo ""
		return 0
	fi

	curr_date=$(date '+%Y-%m-%d %H:%M:%S')
	current_branch=$(git -C "$dir" rev-parse --abbrev-ref HEAD)

	if git -C "$dir" add -A; then
		if git -C "$dir" commit -m "vault backup: $curr_date"; then
			if git -C "$dir" push origin "$current_branch"; then
				echo "Backup for '$dir' committed and updated successfully at $curr_date"
				return 0
			else
				echo "Error: Failed to push to remote for '$dir'"
			fi
		else
			echo "Error: Failed to commit changes for '$dir'"
		fi
	else
		echo "Error: Failed to stage changes for '$dir'"
	fi

	return 1
}

for dir in $DIRS; do
	echo "Processing $dir"
	backup_with_git "$dir" || echo "Backup failed for $dir"
done
