# bin/crontab/backup-music
#!/bin/env sh

set -x
set -e

BACKUP_SOURCE="$HOME/Music"
BACKUP_DESTINATION="/run/media/bryant/tiagovault/backups"
BACKUP_DIR_NAME="music"

backup() {
	backup_folder="$BACKUP_DESTINATION/$BACKUP_DIR_NAME"

	mkdir -p "$backup_folder"

	if rsync -a --delete "$BACKUP_SOURCE/" "$backup_folder/"; then
		sync
		notify-send "Backup music" "finalized successfully"
		sound-alert
	else
		notify-send "Backup music" "failed"
		return 1
	fi
}

backup
